# pdu_mqtt/rootfs/etc/services.d/pdu-mqtt/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# LogiLink & Intellinet PDU MQTT Bridge
# Handles service shutdown
# ==============================================================================

bashio::log.info "PDU MQTT Bridge service stopped"

# Exit with code 0 to prevent restart
exit 0 